# root/profile.d/init
# just the useful stuff from /etc/profile,
# without the overriding of $PS1 rubbish,
# and with the addition of $dir/init.sh bundles

if [ -d /etc/profile.d ]; then
	for i in /etc/profile.d/*.sh; do
		if [ -r $i ]; then
			. $i
		fi
	done
	for i in /etc/profile.d/*; do
		if [ -d $i -a -r $i/init.sh ]; then
			. $i/init.sh
		fi
	done
	unset i
fi
